# systemd/nextjs.service
[Unit]
Description=Next.js SSR app
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/srv/nextjs
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
