# systemd/nextjs.service
[Service]
Type=simple
User=ubuntu
WorkingDirectory=/srv/nextjs
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=0.0.0.0
ExecStart=/usr/bin/env PORT=$PORT HOST=$HOST /usr/bin/npm run start
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
